# configs/deploy/systemd/r-ems-ui.service
# ---
# ems_section: "12-gui-setup-wizard"
# ems_subsection: "systemd-ui"
# ems_type: "config"
# ems_scope: "operations"
# ems_description: "Systemd unit to host the R-EMS UI service."
# ems_version: "v0.0.0-prealpha"
# ems_owner: "tbd"
# ---
[Unit]
Description=R-EMS Log Viewer UI
After=multi-user.target
ConditionPathExists=/usr/local/bin/r-ems-ui

[Service]
Type=simple
User=rems
Group=rems
EnvironmentFile=-/etc/r-ems/r-ems-ui.env
ExecStart=/usr/local/bin/r-ems-ui --dir ${R_EMS_UI_LOG_DIR:-/var/log/r-ems}
WorkingDirectory=/var/lib/r-ems
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target