# configs/deploy/systemd/r-emsd.service
# ---
# ems_section: "01-core-functionality"
# ems_subsection: "systemd-daemon"
# ems_type: "config"
# ems_scope: "operations"
# ems_description: "Systemd unit for the primary R-EMS daemon."
# ems_version: "v0.0.0-prealpha"
# ems_owner: "tbd"
# ---
[Unit]
Description=R-EMS Core Daemon
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rems
Group=rems
EnvironmentFile=-/etc/r-ems/r-emsd.env
ExecStart=/usr/local/bin/r-emsd --config ${R_EMS_CONFIG_PATH:-/etc/r-ems/config.toml}
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target